# trading bot 1.0/com.lucra.tradingbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.lucra.tradingbot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/liamblackshaw-brown/code/Lucra 1.1/trading bot 1.0/.venv/bin/python3</string>
    <string>/Users/liamblackshaw-brown/code/Lucra 1.1/trading bot 1.0/main.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/liamblackshaw-brown/code/Lucra 1.1/trading bot 1.0</string>

  <!-- Restart automatically if the process exits -->
  <key>KeepAlive</key>
  <true/>

  <!-- Throttle restarts to avoid spin-loop -->
  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/Users/liamblackshaw-brown/code/Lucra 1.1/trading bot 1.0/logs/stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/liamblackshaw-brown/code/Lucra 1.1/trading bot 1.0/logs/stderr.log</string>

  <!-- Run at load (i.e. at login) -->
  <key>RunAtLoad</key>
  <true/>

</dict>
</plist>
